# server/awg-fw-watchdog.service
[Unit]
Description=Ensure the AmneziaWG hardened firewall is loaded
Documentation=file:///etc/nftables.conf
After=nftables.service network-pre.target
Wants=nftables.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/awg-fw-ensure.sh
# Least privilege: only needs to read the config and program nftables.
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths=/proc/sys
CapabilityBoundingSet=CAP_NET_ADMIN
AmbientCapabilities=CAP_NET_ADMIN
NoNewPrivileges=yes
